# Arranque máquina virtual automáticamente/Con_otro_usuario/virtualbox/lanzarMaquina.service
[Unit]
Description=Lanzadera de maquina virtual
Before=vboxdrv.service

[Service]
Type=simple
ExecStart=/usr/bin/xinit "/opt/servicioMaquina/lanzarMaquina.sh" -- :1
ExecReload=/usr/bin/xinit "/opt/servicioMaquina/lanzarMaquina.sh" -- :1
Restart=always

[Install]
WantedBy=vboxdrv.service
WantedBy=graphical.target